# deploy/systemd/v5-auto-risk-eval.service
[Unit]
Description=V5 Auto Risk Guard Evaluation

[Service]
Type=oneshot
WorkingDirectory=/home/ubuntu/clawd/v5-prod
Environment=PYTHONUNBUFFERED=1
Environment=PYTHONPATH=/home/ubuntu/clawd/v5-prod
Environment=PATH=/home/ubuntu/clawd/v5-prod/.venv/bin:/usr/local/bin:/usr/bin:/bin
ExecStart=/home/ubuntu/clawd/v5-prod/.venv/bin/python /home/ubuntu/clawd/v5-prod/scripts/auto_risk_eval.py
