# deploy/systemd/v5-cost-rollup.service
[Unit]
Description=V5 cost rollup (yesterday UTC)
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/home/admin/clawd/v5-trading-bot
User=admin
Group=admin
ExecStart=/usr/bin/flock -n /tmp/v5_cost_rollup.lock -c "/bin/bash -lc 'python3 scripts/rollup_costs.py --day $(date -u -d \"yesterday\" +%%Y%%m%%d)'"
