# deploy/systemd/v5-cost-rollup-real.user.service
[Unit]
Description=V5: rollup real (okx_fill) cost stats + anomaly check

[Service]
Type=oneshot
WorkingDirectory=/home/ubuntu/clawd/v5-prod
Environment=PYTHONUNBUFFERED=1
ExecStart=/bin/bash -lc 'cd /home/ubuntu/clawd/v5-prod && /home/ubuntu/clawd/v5-prod/.venv/bin/python scripts/rollup_costs.py --day $(date -u +%%Y%%m%%d) --base_dir reports/cost_events --out_dir reports/cost_stats_real --source okx_fill --check_anomaly'
