# deploy/systemd/v5-daily-ml-training.service
[Unit]
Description=V5 Daily ML Model Training
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/home/admin/clawd/v5-prod
Environment=PYTHONUNBUFFERED=1
Environment=PYTHONPATH=/home/admin/clawd/v5-prod
Environment=V5_WORKSPACE=/home/admin/clawd/v5-prod
Environment=V5_CONFIG=configs/live_prod.yaml
Environment=V5_ML_MIN_SAMPLES=200
Environment=V5_ML_MIN_VALID_IC=0.00
Environment=V5_ML_MAX_IC_GAP=0.25
Environment=V5_ML_MIN_CV_MEAN_IC=0.01
Environment=V5_ML_MAX_CV_STD=0.15
Environment=V5_ML_TARGET_MODE=forward_edge_rank
Environment=V5_ML_CANDIDATES=ridge
Environment=V5_ML_RIDGE_ALPHA=50
Environment=V5_ML_MIN_SYMBOL_SAMPLES=48
Environment=V5_ML_MIN_GROUP_SIZE=2
Environment=V5_ML_MIN_GROUP_COVERAGE_RATIO=0.9
Environment=V5_ML_FEATURE_SELECTOR=stable
Environment=V5_ML_ROLLING_WINDOW_DAYS=10
Environment=V5_ML_RECENCY_HALFLIFE_DAYS=5
Environment=V5_ML_RECENCY_MAX_WEIGHT=3.0
ExecStartPre=-/bin/bash -lc 'cd /home/admin/clawd/v5-prod && source .venv/bin/activate && python scripts/backfill_ml_multihorizon_labels.py'
ExecStart=/bin/bash -lc 'cd /home/admin/clawd/v5-prod && source .venv/bin/activate && python scripts/daily_ml_training.py'

[Install]
WantedBy=default.target
