# deploy/systemd/v5-daily-ml-training.timer
[Unit]
Description=Run V5 Daily ML Training at 00:30 (research-only; disabled in live_prod)
# ML training is kept for offline research and is not active in live_prod.
# Create /etc/v5-enable-ml-research-timers only for an explicit research host opt-in.
ConditionPathExists=/etc/v5-enable-ml-research-timers

[Timer]
OnCalendar=*-*-* 00:30:00
Persistent=true
AccuracySec=60s
RandomizedDelaySec=300s
Unit=v5-daily-ml-training.service

[Install]
WantedBy=timers.target
