# deploy/systemd/v5-export-bundle.timer
[Unit]
Description=Export V5 telemetry bundle every 10 minutes

[Timer]
OnBootSec=3min
OnUnitActiveSec=10min
Persistent=true

[Install]
WantedBy=timers.target
