# deploy/systemd/v5-export-bundle.service
[Unit]
Description=V5 export sanitized telemetry bundle for quant-lab
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=ubuntu
WorkingDirectory=/home/ubuntu/clawd/v5-prod
ExecStart=/home/ubuntu/clawd/v5-prod/.venv/bin/python scripts/export_v5_bundle.py --reports-dir reports --out-dir /var/lib/v5/exports/bundles --window-hours 72
NoNewPrivileges=true
PrivateTmp=true
