# deploy/systemd/v5-hourly.timer
[Unit]
Description=Run V5 dry-run hourly (UTC window id)

[Timer]
OnCalendar=hourly
Persistent=true
Unit=v5-hourly.service

[Install]
WantedBy=timers.target
