# deploy/systemd/v5-live-followup-bundle-export.service
[Unit]
Description=Generate sanitized V5 live follow-up bundle for quant-lab
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/home/ubuntu/clawd/v5-prod
Environment=V5_WORKSPACE=/home/ubuntu/clawd/v5-prod
Environment=V5_LIVE_FOLLOWUP_BUNDLE_DIR=/var/lib/v5/exports/bundles
Environment=V5_BUNDLE_SOURCE_USER=ubuntu
Environment=V5_DASHBOARD_BUNDLE_USER=ubuntu
ExecStart=/bin/bash /home/ubuntu/clawd/v5-prod/scripts/install_live_followup_bundle.sh
NoNewPrivileges=true
PrivateTmp=true
