# deploy/systemd/v5-live-followup-bundle-export.timer
[Unit]
Description=Generate V5 live follow-up bundle every 10 minutes for quant-lab ingest

[Timer]
OnBootSec=3min
OnUnitActiveSec=10min
Persistent=true
AccuracySec=30s
RandomizedDelaySec=2m
Unit=v5-live-followup-bundle-export.service

[Install]
WantedBy=timers.target
