# deploy/systemd/v5-prod.user.service
[Unit]
Description=V5 production LIVE hourly window (OKX spot)
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/home/admin/clawd/v5-trading-bot
Environment=PYTHONUNBUFFERED=1
Environment=V5_CONFIG=configs/live_prod.yaml
Environment=V5_DATA_PROVIDER=okx
Environment=V5_LIVE_ARM=YES
ExecStart=/usr/bin/flock -n /tmp/v5_live_prod_hourly.lock -c "/bin/bash -lc '/home/admin/clawd/v5-trading-bot/scripts/run_hourly_live_window.sh'"
