# deploy/systemd/v5-quant-lab-bundle-sync.timer
[Unit]
Description=Sync latest quant-lab expert pack every 10 minutes

[Timer]
OnBootSec=3min
OnUnitActiveSec=10min
Persistent=true

[Install]
WantedBy=timers.target
