# deploy/systemd/v5-reconcile.timer
[Unit]
Description=Run V5 reconcile every 5 minutes

[Timer]
# Calendar timer so Persistent=true can catch up missed runs after downtime.
OnCalendar=*:0/5
Persistent=true
AccuracySec=10s
RandomizedDelaySec=10s
Unit=v5-reconcile.service

[Install]
WantedBy=timers.target
