# deploy/systemd/v5-reconcile.timer
[Unit]
Description=Run V5 reconcile every 5 minutes

[Timer]
OnUnitActiveSec=300
Persistent=true
AccuracySec=10s
RandomizedDelaySec=10s
Unit=v5-reconcile.service

[Install]
WantedBy=timers.target
