# deploy/systemd/v5-reconcile.service
[Unit]
Description=V5 reconcile (OKX balance vs local stores)
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/home/ubuntu/clawd/v5-prod
Environment=PYTHONUNBUFFERED=1
SyslogIdentifier=v5-reconcile
ExecStart=/usr/bin/flock -n /tmp/v5_reconcile.lock -c "/bin/bash -lc 'PYTHONPATH=. /home/ubuntu/clawd/v5-prod/.venv/bin/python scripts/reconcile_guard_once.py --config configs/live_prod.yaml --env .env'"
Restart=on-failure
RestartSec=10
