# deploy/systemd/v5-reconcile.service
[Unit]
Description=V5 reconcile (OKX balance vs local stores)
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/home/admin/clawd/v5-trading-bot
User=admin
Group=admin
Environment=PYTHONUNBUFFERED=1
SyslogIdentifier=v5-reconcile
ExecStart=/usr/bin/flock -n /tmp/v5_reconcile.lock -c "/bin/bash -lc '/home/admin/clawd/v5-trading-bot/.venv/bin/python scripts/reconcile_guard_once.py --out reports/reconcile_status.json'"
Restart=on-failure
RestartSec=10
