# deploy/systemd/v5-spread-rollup.user.service
[Unit]
Description=V5 spread rollup (yesterday UTC)

[Service]
Type=oneshot
WorkingDirectory=/home/admin/clawd/v5-trading-bot
Environment=PYTHONUNBUFFERED=1
ExecStart=/bin/bash -lc 'cd /home/admin/clawd/v5-trading-bot && /home/admin/clawd/v5-trading-bot/.venv/bin/python scripts/rollup_spreads.py --day $(date -u -d "yesterday" +%%Y%%m%%d)'
