# deploy/systemd/v5-sentiment-collect.service
[Unit]
Description=V5 Sentiment Data Collection

[Service]
Type=oneshot
WorkingDirectory=/home/admin/clawd/v5-trading-bot
Environment=PYTHONUNBUFFERED=1
ExecStart=/bin/bash -lc 'cd /home/admin/clawd/v5-trading-bot && .venv/bin/python scripts/collect_funding_sentiment.py && .venv/bin/python scripts/collect_rss_sentiment.py'
