# deploy/systemd/v5-shadow-regime.user.timer
[Unit]
Description=Run V5 shadow regime-aware weights hourly

[Timer]
OnCalendar=hourly
Persistent=true
AccuracySec=30s
RandomizedDelaySec=25s
Unit=v5-shadow-regime.user.service

[Install]
WantedBy=timers.target
