# deploy/systemd/v5-spread-rollup.service
[Unit]
Description=V5 spread rollup (yesterday UTC)
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/home/admin/clawd/v5-trading-bot
User=admin
Group=admin
Environment=PYTHONUNBUFFERED=1
SyslogIdentifier=v5-spread-rollup
ExecStart=/usr/bin/flock -n /tmp/v5_spread_rollup.lock -c "/bin/bash -lc '/home/admin/clawd/v5-trading-bot/.venv/bin/python scripts/rollup_spreads.py --day $(date -u -d \"yesterday\" +%%Y%%m%%d)'"
