# deploy/systemd/v5-shadow-regime.user.service
[Unit]
Description=V5 shadow (dry-run) regime-aware weights hourly evaluation
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/home/admin/clawd/v5-trading-bot
Environment=PYTHONUNBUFFERED=1
Environment=V5_CONFIG=configs/live_20u_shadow_regime.yaml
Environment=V5_DATA_PROVIDER=okx
Environment=V5_LIVE_ARM=NO
Environment=V5_DISABLE_TOPLEVEL_ARTIFACTS=1
ExecStart=/usr/bin/flock -n /tmp/v5_shadow_regime_hourly.lock -c "/bin/bash -lc '/home/admin/clawd/v5-trading-bot/scripts/run_hourly_live_window.sh'"
