# deploy/systemd/v5-trade-monitor.service
[Unit]
Description=V5 Trading Monitor - Check trade health and send alerts
After=network.target

[Service]
Type=oneshot
WorkingDirectory=/home/ubuntu/clawd/v5-prod
Environment="PATH=/home/ubuntu/clawd/v5-prod/.venv/bin:/usr/local/bin:/usr/bin:/bin"
Environment="PYTHONPATH=/home/ubuntu/clawd/v5-prod"
ExecStart=/home/ubuntu/clawd/v5-prod/.venv/bin/python /home/ubuntu/clawd/v5-prod/scripts/v5_trade_monitor.py --silent
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
