# deploy/systemd/v5-trade-monitor.timer
[Unit]
Description=Run V5 trade monitor every hour

[Timer]
OnCalendar=*-*-* *:07:00
Persistent=true
AccuracySec=30s
Unit=v5-trade-monitor.service

[Install]
WantedBy=timers.target
