# archive/20260313-research-cleanup/deploy/systemd/v5-shadow-tuned-xgboost.user.timer
[Unit]
Description=Run V5 shadow tuned XGBoost hourly

[Timer]
OnCalendar=hourly
Persistent=true
AccuracySec=30s
RandomizedDelaySec=35s
Unit=v5-shadow-tuned-xgboost.user.service

[Install]
WantedBy=timers.target
